# system/systemd/lva-audio-watchdog.timer
[Unit]
Description=Run LVA Audio Watchdog every 2 minutes

[Timer]
OnBootSec=120
OnUnitActiveSec=120

[Install]
WantedBy=timers.target
